# 5001/testapp-hello-everyone.conf
[program:testapp]
command =  python /opt/bcgh-app/app-hello-everyone.py 
