# testapp.conf
[program:testapp]
command =  python /opt/bcgh-app/app.py 
